# systemd/escpos-bridge.service
[Unit]
Description=EposNow ESC/POS print stream interceptor (receipt + KP)
Documentation=https://github.com/jmcbride4882/escpos-bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/escpos-bridge/config.env
WorkingDirectory=/opt/escpos-bridge
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
User=nobody
Group=nogroup
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/escpos-bridge
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true

[Install]
WantedBy=multi-user.target
