# systemd/escpos-bridge.service
[Unit]
Description=EposNow ESC/POS print stream interceptor (receipt + KP)
Documentation=https://github.com/jmcbride4882/escpos-bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/escpos-bridge/config.env
WorkingDirectory=/opt/escpos-bridge
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Runs as a dedicated 'escpos' user so the web GUI can write config.env
# and call `sudo systemctl restart escpos-bridge` (sudoers entry created
# by install.sh). install.sh creates the user.
User=escpos
Group=escpos

# Hardening (relaxed enough that the web GUI can edit /etc/escpos-bridge/)
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/escpos-bridge /etc/escpos-bridge
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
LockPersonality=true

[Install]
WantedBy=multi-user.target
